# services.d/sockd/run
#!/bin/sh

while [ "$(/sbin/ip link show | /usr/bin/cut -d' ' -f2 | /bin/grep -o tun)" == "" ];
do
  /bin/sleep 1
done

exec /usr/sbin/sockd
